# web-integrity-monitor.timer
[Unit]
Description=Web Integrity Monitor Timer (Every 5 minutes)
Requires=web-integrity-monitor.service

[Timer]
# Start 2 minutes after boot
OnBootSec=2min

# Run every 5 minutes after the last activation
OnUnitActiveSec=5min

# Accuracy of timer (lower = more precise, higher = more efficient)
AccuracySec=1s

[Install]
WantedBy=timers.target
